# ccaa_dir/services/ccaa_web.service
[Unit]
Description=ccaa_web
After=network.target

[Service]
ExecStart=/usr/bin/sudo -u ccaa /usr/sbin/ccaa_web

[Install]
WantedBy=multi-user.target